# utils/keepalive.service
[Service]
ExecStart=/usr/local/bin/.keepalive.sh
User=admin

[Install]
WantedBy=multi-user.target
